# youtube-export/crontab
MAILTO = toby-admin+crontab@khanacademy.org

1 0 * * * $HOME/youtube-export/cron/convert.sh || echo "Failed to run youtube-export; see $HOME/youtube-export/logs"
